# systemd/zimmer-user-panel.service
[Unit]
Description=Zimmer Platform User Panel
After=network.target zimmer-backend.service
Requires=zimmer-backend.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/zimmer/zimmer_user_panel
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npm start
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zimmer-user-panel

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
